# build/rhoconnect-push.conf
description "rhoconnect-push server"
author      "A. Babichev"
 
start on startup
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5
 
script
    export PATH=/opt/rhoconnect/bin:$PATH 
    exec /opt/rhoconnect/bin/rhoconnect-push -c /opt/rhoconnect/etc/config.json >> /var/log/rhoconnect-push.log 2>&1
end script
 
post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "rhoconnect-push has started!"
end script




